# local_rules.xml
<!-- ============================================================
     SOC SIEM Home Lab - Custom Wazuh Detection Rules
     Location: /var/ossec/etc/rules/local_rules.xml
     Author: SOC Home Lab Project
     MITRE ATT&CK Mapped
     ============================================================ -->

<group name="local,soc-homelab,">

  <!-- ──────────────────────────────────────────────
       RULE 1: Brute Force Detection
       MITRE: T1110 - Brute Force
       Fires when: 5+ failed logins in 60 seconds
       Level 10 = HIGH severity
       ────────────────────────────────────────────── -->
  <rule id="100001" level="10" frequency="5" timeframe="60">
    <if_matched_sid>60122</if_matched_sid>
    <description>Brute force attack detected - multiple failed logins from same source</description>
    <mitre>
      <id>T1110</id>
    </mitre>
  </rule>

  <!-- ──────────────────────────────────────────────
       RULE 2: Port Scan Detection
       MITRE: T1046 - Network Service Discovery
       Fires when: High volume Sysmon network events
       Level 8 = MEDIUM severity
       ────────────────────────────────────────────── -->
  <rule id="100002" level="8">
    <if_group>sysmon</if_group>
    <field name="win.system.eventID">3</field>
    <description>Possible port scan detected - high volume network connections from single source</description>
    <mitre>
      <id>T1046</id>
    </mitre>
  </rule>

  <!-- ──────────────────────────────────────────────
       RULE 3: PowerShell Encoded Command
       MITRE: T1059.001 - PowerShell
       Fires when: PowerShell run with -enc flag
       Level 12 = CRITICAL severity
       ────────────────────────────────────────────── -->
  <rule id="100003" level="12">
    <if_group>windows</if_group>
    <field name="win.eventdata.commandLine" type="pcre2">-[Ee][nN][cC]</field>
    <description>PowerShell encoded command detected - possible malware or living-off-the-land attack</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
  </rule>

  <!-- ──────────────────────────────────────────────
       RULE 4: LSASS Access - Credential Dumping
       MITRE: T1003.001 - LSASS Memory
       Fires when: Any process accesses lsass.exe
       Level 15 = MAXIMUM severity (drop everything)
       ────────────────────────────────────────────── -->
  <rule id="100004" level="15">
    <if_group>sysmon</if_group>
    <field name="win.eventdata.targetImage" type="pcre2">lsass\.exe</field>
    <field name="win.system.eventID">10</field>
    <description>LSASS process access detected - credential dumping attempt (Mimikatz/similar)</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
  </rule>

  <!-- ──────────────────────────────────────────────
       RULE 5: Suspicious Process - CMD from Office App
       MITRE: T1566 - Phishing
       Fires when: cmd.exe spawned by Office application
       Level 12 = CRITICAL severity
       ────────────────────────────────────────────── -->
  <rule id="100005" level="12">
    <if_group>sysmon</if_group>
    <field name="win.system.eventID">1</field>
    <field name="win.eventdata.image" type="pcre2">cmd\.exe</field>
    <field name="win.eventdata.parentImage" type="pcre2">WINWORD\.EXE|EXCEL\.EXE|POWERPNT\.EXE</field>
    <description>Suspicious process - CMD spawned by Microsoft Office application (possible macro attack)</description>
    <mitre>
      <id>T1566</id>
    </mitre>
  </rule>

  <!-- ──────────────────────────────────────────────
       RULE 6: New Admin User Created
       MITRE: T1136 - Create Account
       Fires when: New user added to Administrators
       Level 12 = CRITICAL severity
       ────────────────────────────────────────────── -->
  <rule id="100006" level="12">
    <if_group>windows</if_group>
    <field name="win.system.eventID">4720</field>
    <description>New user account created - possible unauthorized account creation</description>
    <mitre>
      <id>T1136.001</id>
    </mitre>
  </rule>

  <!-- ──────────────────────────────────────────────
       RULE 7: Scheduled Task Created
       MITRE: T1053 - Scheduled Task/Job
       Fires when: New scheduled task registered
       Level 10 = HIGH severity
       ────────────────────────────────────────────── -->
  <rule id="100007" level="10">
    <if_group>windows</if_group>
    <field name="win.system.eventID">4698</field>
    <description>Scheduled task created - possible persistence mechanism</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
  </rule>

  <!-- ──────────────────────────────────────────────
       RULE 8: Windows Defender Disabled
       MITRE: T1562.001 - Disable Security Tools
       Fires when: Defender tampered with
       Level 15 = MAXIMUM severity
       ────────────────────────────────────────────── -->
  <rule id="100008" level="15">
    <if_group>windows</if_group>
    <field name="win.eventdata.commandLine" type="pcre2">Set-MpPreference.*DisableRealtimeMonitoring</field>
    <description>Windows Defender disabled via PowerShell - active defense evasion detected</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
  </rule>

</group>
